# dynamic-layers/meta-fus/recipes-example/slint-demos/slint-demos/slint-demos.service
[Unit]
Description=Startup script to launch Slint demo(s) on startup

[Service]
ExecStart=/usr/bin/energy-monitor

[Install]
WantedBy=multi-user.target