# dynamic-layers/meta-ti/recipes-example/slint-demos/slint-demos/slint-demos.service
[Unit]
Description=Startup script to launch Slint demo(s) on startup
# tidss (the display controller) and its HDMI bridge (sii902x) are loadable
# kernel modules that only probe a few seconds into boot, so /dev/dri/card0
# doesn't exist yet when systemd would otherwise start us -- the demo then fails
# to open the DRM device. Order after the card's device unit, which systemd
# materialises once the DRM device is tagged for it (see 99-drm-device-unit.rules);
# seatd (if present) brokers DRM/input access.
Wants=dev-dri-card0.device
After=dev-dri-card0.device seatd.service

[Service]
ExecStart=/usr/bin/home-automation
# Backstop: if some other early-boot hiccup (e.g. a slightly later display-card
# probe on the GPU boards, which enumerate two DRM nodes) still trips the first
# start, retry rather than leaving the screen dark.
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
